# configs/supervisor.conf
[program:bottle.svetlyak.ru]
command = /home/art/www/bottle.svetlyak.ru/env/bin/uwsgi /home/art/www/bottle.svetlyak.ru/configs/uwsgi.ini
